# myproject.service
[Unit]
Description=uWSGI instance to serve myproject
After=network.target

[Service]
User=k
Group=www-data
WorkingDirectory=/home/k/myproject
Environment="PATH=/home/k/myproject/myprojectenv/bin"
ExecStart=/home/k/myproject/myprojectenv/bin/uwsgi --ini myproject.ini

[Install]
WantedBy=multi-user.target
