# myproject.service
[Unit]
Description=uWSGI instance to serve myproject
After=network.target

[Service]
User=deployer
Group=www-data
WorkingDirectory=/var/myproject
Environment="PATH=/var/myproject/myprojectenv/bin"
ExecStart=/var/myproject/myprojectenv/bin/uwsgi --ini myproject.ini

[Install]
WantedBy=multi-user.target
